# files/etc/uci-defaults/99-h5000m-defaults
#!/bin/sh

logger -t h5000m-defaults "设置默认时区和本地化"
uci -q set system.@system[0].zonename='Asia/Shanghai'
uci -q set system.@system[0].timezone='CST-8'
uci -q get system.ntp >/dev/null 2>&1 || uci -q set system.ntp='timeserver'
uci -q set system.ntp.enabled='1'
uci -q delete system.ntp.server
uci -q add_list system.ntp.server='ntp.aliyun.com'
uci -q add_list system.ntp.server='cn.pool.ntp.org'
uci -q add_list system.ntp.server='time1.cloud.tencent.com'
uci -q add_list system.ntp.server='time.ustc.edu.cn'
uci -q commit system

if [ -f /etc/config/luci ]; then
	uci -q get luci.main >/dev/null 2>&1 || uci -q set luci.main='core'
	uci -q set luci.main.lang='zh_Hans'
	uci -q commit luci
fi

read_uboot_env() {
	local key="$1"
	[ -e /dev/mmcblk0p1 ] || return 1
	tr '\0' '\n' < /dev/mmcblk0p1 | sed -n "s/^${key}=//p" | head -n 1
}

valid_mac() {
	echo "$1" | grep -Eq '^[0-9a-fA-F]{2}(:[0-9a-fA-F]{2}){5}$'
}

logger -t h5000m-defaults "设置默认 LAN 地址为 192.168.10.1"
uci -q set network.lan.proto='static'
uci -q set network.lan.ipaddr='192.168.10.1'
uci -q set network.lan.netmask='255.255.255.0'

logger -t h5000m-defaults "从 U-Boot 环境写入默认 MAC 地址"
lan_mac="$(read_uboot_env ethaddr || true)"
wan_mac="$(read_uboot_env eth1addr || true)"
if valid_mac "${lan_mac}" && valid_mac "${wan_mac}"; then
	uci -q set network.eth0_dev='device'
	uci -q set network.eth0_dev.name='eth0'
	uci -q set network.eth0_dev.macaddr="${lan_mac}"
	uci -q set network.eth1_dev='device'
	uci -q set network.eth1_dev.name='eth1'
	uci -q set network.eth1_dev.macaddr="${wan_mac}"
	uci -q set network.wan.macaddr="${wan_mac}"
	uci -q commit network
else
	logger -t h5000m-defaults "U-Boot MAC 地址无效或不存在，保留系统默认 MAC"
fi

logger -t h5000m-defaults "写入 H5000M QModem 默认槽位"
if [ -f /etc/config/qmodem ]; then
	uci -q set qmodem.usb0='modem-slot'
	uci -q set qmodem.usb0.slot='2-1'
	uci -q set qmodem.usb0.type='usb'
	uci -q set qmodem.usb0.alias='USB'
	uci -q set qmodem.main.enable_pcie_scan='0'
	uci -q set qmodem.main.enable_usb_scan='0'
	uci -q set qmodem.main.try_preset_pcie='0'
	uci -q set qmodem.main.try_preset_usb='1'
	uci -q set qmodem.main.start_delay='0'
	uci -q commit qmodem
fi

logger -t h5000m-defaults "优化 WAN 与 5G 备用出口"
if [ -x /usr/sbin/h5000m-netmode ]; then
	/usr/sbin/h5000m-netmode >/dev/null 2>&1 || true
fi

logger -t h5000m-defaults "设置 root 默认密码"
if command -v chpasswd >/dev/null 2>&1; then
	echo 'root:admin' | chpasswd
else
	(printf 'admin\nadmin\n' | passwd root) >/dev/null 2>&1 || true
fi

logger -t h5000m-defaults "设置默认 WiFi 密码"
if ! uci -q show wireless | grep -q '=wifi-iface'; then
	wifi config >/dev/null 2>&1 || true
fi

for radio in $(uci -q show wireless | sed -n 's/^\(wireless\.[^.]*\)=wifi-device$/\1/p'); do
	uci -q set "${radio}.disabled=0"
	uci -q set "${radio}.country=CN"
done

for section in $(uci -q show wireless | sed -n 's/^\(wireless\.[^.]*\)=wifi-iface$/\1/p'); do
	uci -q set "${section}.disabled=0"
	uci -q set "${section}.ssid=H5000M"
	uci -q set "${section}.encryption=psk2"
	uci -q set "${section}.key=1234567890"
done
uci -q commit wireless

logger -t h5000m-defaults "清理固件内的 QModem 和 video apk 软件源"
if [ -f /etc/apk/repositories.d/distfeeds.list ]; then
	sed -i '/\/qmodem\//d;/ qmodem[ _-]/d;/ qmodem$/d;/\/small_package\//d;/ small_package[ _-]/d;/ small_package$/d;/\/video\//d' /etc/apk/repositories.d/distfeeds.list
fi

exit 0
