# files/etc/uci-defaults/99-h5000m-defaults
#!/bin/sh

[ -f /lib/functions/system.sh ] && . /lib/functions/system.sh

logger -t h5000m-defaults "Set default timezone and locale defaults"
uci -q set system.@system[0].zonename='Asia/Shanghai'
uci -q set system.@system[0].timezone='CST-8'
uci -q get system.ntp >/dev/null 2>&1 || uci -q set system.ntp='timeserver'
uci -q set system.ntp.enabled='1'
uci -q delete system.ntp.server
uci -q add_list system.ntp.server='ntp.aliyun.com'
uci -q add_list system.ntp.server='cn.pool.ntp.org'
uci -q add_list system.ntp.server='time1.cloud.tencent.com'
uci -q add_list system.ntp.server='time.ustc.edu.cn'
uci -q commit system

if [ -f /etc/config/luci ]; then
	uci -q get luci.main >/dev/null 2>&1 || uci -q set luci.main='core'
	uci -q set luci.main.lang='auto'
	uci -q commit luci
fi

if [ -f /etc/config/upnpd ]; then
	logger -t h5000m-defaults "Enable UPnP and NAT-PMP by default"
	uci -q set upnpd.config.enabled='1'
	uci -q set upnpd.config.enable_natpmp='1'
	uci -q set upnpd.config.enable_upnp='1'
	uci -q commit upnpd
fi

read_uboot_env() {
	local key="$1"
	[ -e /dev/mmcblk0p1 ] || return 1
	tr '\0' '\n' < /dev/mmcblk0p1 | sed -n "s/^${key}=//p" | head -n 1
}

valid_mac() {
	echo "$1" | grep -Eq '^[0-9a-fA-F]{2}(:[0-9a-fA-F]{2}){5}$'
}

mac_add() {
	local mac="$1"
	local add="$2"
	local oldifs b1 b2 b3 b4 b5 b6 value

	oldifs="${IFS}"
	IFS=:
	set -- ${mac}
	IFS="${oldifs}"
	b1="$1"
	b2="$2"
	b3="$3"
	b4="$4"
	b5="$5"
	b6="$6"
	value=$(( (0x${b4} << 16) + (0x${b5} << 8) + 0x${b6} + add ))
	value=$(( value & 0xffffff ))
	printf '%02x:%02x:%02x:%02x:%02x:%02x' $(( 0x${b1} )) $(( 0x${b2} )) $(( 0x${b3} )) \
		$(( (value >> 16) & 0xff )) $(( (value >> 8) & 0xff )) $(( value & 0xff ))
}

get_base_mac() {
	local base

	base="$(read_uboot_env ethaddr || true)"
	if valid_mac "${base}"; then
		echo "${base}"
		return 0
	fi

	if type macaddr_generate_from_mmc_cid >/dev/null 2>&1; then
		base="$(macaddr_generate_from_mmc_cid mmcblk0 2>/dev/null || true)"
		if valid_mac "${base}"; then
			logger -t h5000m-defaults "U-Boot ethaddr is invalid, using stable MAC derived from eMMC CID"
			echo "${base}"
			return 0
		fi
	fi

	return 1
}

logger -t h5000m-defaults "Set default LAN address to 10.0.0.1"
uci -q set network.lan.proto='static'
uci -q set network.lan.ipaddr='10.0.0.1'
uci -q set network.lan.netmask='255.255.255.0'

logger -t h5000m-defaults "Apply official H5000M Ethernet layout: eth1 LAN, eth0 WAN"
for section in $(uci -q show network | sed -n 's/^\(network\.[^=]*\)=device$/\1/p'); do
	if [ "$(uci -q get "${section}.name")" = "br-lan" ]; then
		uci -q set "${section}.ports=eth1"
	fi
done
uci -q set network.lan.device='br-lan'
uci -q set network.wan.device='eth0'
uci -q set network.wan6.device='eth0'

logger -t h5000m-defaults "Apply MAC addresses derived from U-Boot ethaddr"
lan_mac="$(get_base_mac || true)"
if valid_mac "${lan_mac}"; then
	wan_mac="$(mac_add "${lan_mac}" 1)"
	wifi24_mac="$(mac_add "${lan_mac}" 2)"
	wifi5_mac="$(mac_add "${lan_mac}" 3)"

	uci -q set network.eth0_dev='device'
	uci -q set network.eth0_dev.name='eth0'
	uci -q set network.eth0_dev.macaddr="${wan_mac}"
	uci -q set network.eth1_dev='device'
	uci -q set network.eth1_dev.name='eth1'
	uci -q set network.eth1_dev.macaddr="${lan_mac}"
	uci -q set network.wan.macaddr="${wan_mac}"
	uci -q commit network
else
	logger -t h5000m-defaults "No valid U-Boot or eMMC CID MAC source found, keeping system defaults"
fi

logger -t h5000m-defaults "Write default H5000M QModem slot"
if [ -f /etc/config/qmodem ]; then
	uci -q set qmodem.usb0='modem-slot'
	uci -q set qmodem.usb0.slot='2-1'
	uci -q set qmodem.usb0.type='usb'
	uci -q set qmodem.usb0.alias='USB'
	uci -q set qmodem.main.enable_pcie_scan='0'
	uci -q set qmodem.main.enable_usb_scan='0'
	uci -q set qmodem.main.try_preset_pcie='0'
	uci -q set qmodem.main.try_preset_usb='1'
	uci -q set qmodem.main.start_delay='0'
	uci -q commit qmodem
fi

logger -t h5000m-defaults "Apply WAN and 5G default route policy"
if [ -x /usr/sbin/h5000m-netmode ]; then
	/usr/sbin/h5000m-netmode >/dev/null 2>&1 || true
fi

logger -t h5000m-defaults "Set default root password"
if command -v chpasswd >/dev/null 2>&1; then
	echo 'root:admin' | chpasswd
else
	(printf 'admin\nadmin\n' | passwd root) >/dev/null 2>&1 || true
fi

logger -t h5000m-defaults "Set separated default WiFi SSIDs and password"
if ! uci -q show wireless | grep -q '=wifi-iface'; then
	wifi config >/dev/null 2>&1 || true
fi

for radio in $(uci -q show wireless | sed -n 's/^\(wireless\.[^.]*\)=wifi-device$/\1/p'); do
	band="$(uci -q get "${radio}.band")"
	uci -q set "${radio}.country=CN"
	case "${band}" in
		2g)
			uci -q set "${radio}.disabled=0"
			uci -q set "${radio}.htmode=EHT40"
			;;
		5g)
			uci -q set "${radio}.disabled=0"
			uci -q set "${radio}.htmode=EHT160"
			;;
		*)
			uci -q set "${radio}.disabled=0"
			;;
	esac
done

wifi_index=2
for section in $(uci -q show wireless | sed -n 's/^\(wireless\.[^.]*\)=wifi-iface$/\1/p'); do
	device="$(uci -q get "${section}.device")"
	band="$(uci -q get "wireless.${device}.band")"
	case "${band}" in
		2g)
			uci -q set "${section}.ssid=H5000M-2G"
			;;
		5g|6g)
			uci -q set "${section}.ssid=H5000M-5G"
			;;
		*)
			uci -q set "${section}.ssid=H5000M"
			;;
	esac
	uci -q set "${section}.encryption=sae-mixed"
	uci -q set "${section}.key=77778888"
	uci -q set "${section}.ieee80211w=1"
	uci -q set "${section}.ieee80211k=1"
	uci -q set "${section}.rrm_neighbor_report=1"
	uci -q set "${section}.rrm_beacon_report=1"
	case "${band}" in
		*)
			uci -q set "${section}.disabled=0"
			;;
	esac
	if valid_mac "${lan_mac:-}"; then
		case "${band}" in
			2g)
				uci -q set "${section}.macaddr=${wifi24_mac}"
				;;
			5g|6g)
				uci -q set "${section}.macaddr=${wifi5_mac}"
				;;
			*)
				uci -q set "${section}.macaddr=$(mac_add "${lan_mac}" "${wifi_index}")"
				wifi_index=$((wifi_index + 1))
				;;
		esac
	fi
done
uci -q commit wireless

logger -t h5000m-defaults "Remove build-only apk feed entries"
if [ -f /etc/apk/repositories.d/distfeeds.list ]; then
	sed -i '/\/qmodem\//d;/ qmodem[ _-]/d;/ qmodem$/d;/\/small_package\//d;/ small_package[ _-]/d;/ small_package$/d;/\/video\//d' /etc/apk/repositories.d/distfeeds.list
fi

exit 0
