# packages/at-webserver/files/etc/uci-defaults/91-at-webserver-h5000m
#!/bin/sh

uci -q set at-webserver.config.enabled='1'
uci -q set at-webserver.config.connection_type='SERIAL'
uci -q set at-webserver.config.serial_port='/dev/ttyUSB1'
uci -q set at-webserver.config.serial_port_custom='/dev/ttyUSB1'
uci -q set at-webserver.config.serial_baudrate='115200'
uci -q set at-webserver.config.serial_timeout='10'
uci -q set at-webserver.config.websocket_port='8765'
uci -q set at-webserver.config.websocket_allow_wan='0'
uci -q commit at-webserver

/etc/init.d/at-webserver enable >/dev/null 2>&1 || true

exit 0
